# docker/supervisor/crontab/laravel
* * * * * cd /var/www/html && /usr/local/bin/php artisan schedule:run >> /var/www/html/storage/logs/laravel-schedule.log 2>&1
