# docker/supervisor/supervisord.d/cron.conf
[program:cron]
process_name=%(program_name)s
command=/usr/sbin/crond -f -x sch,proc
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/cron.log
